# gs_8600.conf
# gs_8600.conf

start on startup
stop on shutdown

#respawn indicates it would restart everytime it is crashed of killed,
#in spite of using 'stop'
respawn
respawn limit 9999 0

limit nofile 20480 20480

instance $PORT,$WORLD,$TARGET

script
hostname=$(hostname)
  nohup node /data/nba/nba_game_server/app/game_node $PORT $WORLD $TARGET DEBUG >> /data/nba/log/$WORLD-$hostname-game1.log 2>&1
end script
